# 15-letsencrypt/ansible/roles/node-exporter/templates/node-exporter.service.j2
[Unit]
Description=Node Exporter
After=network-online.target

[Service]
User=node-exporter
Group=node-exporter
Type=simple
ExecStart=/usr/local/bin/node-exporter

[Install]
WantedBy=multi-user.target

